# LumaNet/LumaNode/lumanet.service
[Unit]
Description=LumaNet Distributed LMS Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/LumaNodeV2/server
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /home/pi/LumaNodeV2/server/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
